# gdrive-sync.service
[Unit]
Description=Google Drive media sync service
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
EnvironmentFile=/etc/gdrive2video.env
WorkingDirectory=%h/gdrive2video
ExecStart=/usr/bin/python3 %h/gdrive2video/gdrive_sync.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
