# gdrive-sync.timer
[Unit]
Description=Google Drive media sync timer
Requires=gdrive-sync.service

[Timer]
# Run on boot after 1 minute
OnBootSec=1min
# Run every 5 minutes
OnUnitActiveSec=5min
# Run at a consistent time if missed
Persistent=true

[Install]
WantedBy=timers.target
